# huaweiAktBbo-udev.rules
# Activate HUAWEI E220 modem
#
SUBSYSTEM=="usb", ATTR{idProduct}=="1003", ATTR{idVendor}=="12d1", RUN+="/usr/sbin/huaweiAktBbo"
